# app-slug.service
[Unit]
Description=APP_NAME
After=network.target

[Service]
Type=simple
User=dhughes
WorkingDirectory=/home/dhughes/apps/APP_SLUG
ExecStart=/home/dhughes/apps/APP_SLUG/venv/bin/python app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
